# templates/flaskpgsql.conf
start on runlevel [2345]
stop on runlevel [016]

respawn
setuid nobody
setgid nogroup
exec python /srv/flask-pgsql-demo/app.py
